# .config/cron/root/crontab
0 12 * * * ~/.config/cron/root/daily.sh
