# .config/cron/root/crontab
0 12 * * * ~/.config/cron/daily_root.sh
